# sz.service
[Unit]
Description=SelectZyme app
After=docker.service
Before=

[Service]
User=root
Group=root
Type=oneshot
RemainAfterExit=yes
TimeoutStartSec=600s
ExecStart=LOCAL_SZ_INSTALLATION_SOURCE/sz.sh start
ExecStop=LOCAL_SZ_INSTALLATION_SOURCE/sz.sh stop

[Install]
WantedBy=multi-user.target